# bin/crontab
55 8 * * 1-5 cd /home/dev/src/ctp/CTPService/bin/ && ./startQ.sh
59 8 * * 1-5 cd /home/dev/src/ctp/CTPService/bin/ && ./ctpService start ctp_w
31 11 * * 1-5 cd /home/dev/src/ctp/CTPService/bin/ && ./ctpService stop
0 13 * * 1-5 cd /home/dev/src/ctp/CTPService/bin/ && ./stopQ.sh

25 13 * * 1-5 cd /home/dev/src/ctp/CTPService/bin/ && ./startQ.sh
29 13 * * 1-5 cd /home/dev/src/ctp/CTPService/bin/ && ./ctpService start ctp_w
1 15 * * 1-5 cd /home/dev/src/ctp/CTPService/bin/ && ./ctpService stop
30 16 * * 1-5 cd /home/dev/src/ctp/CTPService/bin/ && ./stopQ.sh

55 20 * * 1-5 cd /home/dev/src/ctp/CTPService/bin/ && ./startQ.sh
59 20 * * 1-5 cd /home/dev/src/ctp/CTPService/bin/ && ./ctpService start ctp_w
31 2 * * 2-6 cd /home/dev/src/ctp/CTPService/bin/ && ./ctpService stop
0 4 * * 2-6 cd /home/dev/src/ctp/CTPService/bin/ && ./stopQ.sh

# 获取手续费定时脚本
5 9 * * 1-5 cd /home/dev/src/ctp/CTPService/bin/ && ./cron/qryRate.py
35 13 * * 1-5 cd /home/dev/src/ctp/CTPService/bin/ && ./cron/qryRate.py
5 21 * * 1-5 cd /home/dev/src/ctp/CTPService/bin/ && ./cron/qryRate.py
